# sshproxy_syncfile.cron
PATH=/bin:/usr/bin
MAILTO=root

* * * * *       sshapp    if [ -x /usr/bin/sshproxy_syncfile ]; then nice -n10 /usr/bin/sshproxy_syncfile; fi

# EOF

